# config/systemd/user/dav-background-sync-one.timer
[Unit]
Description=Sync dav regularly

[Timer]
OnBootSec=60s
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
